# apps/cron/crontab
# Log rotation
@daily /usr/sbin/logrotate /etc/logrotate.d/app


# Sample (every minute)
* * * * * cd /usr/src/app/apps/cron && npm run -s clear-unused-uploads >> /var/log/cron.log 2>&1
